# com.fueltracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Installation als LaunchDaemon (läuft auch ohne Login):
  ───────────────────────────────────────────────────────
  Alten LaunchAgent zuerst entfernen:
    launchctl bootout gui/$(id -u)/com.fueltracker
    rm ~/Library/LaunchAgents/com.fueltracker.plist

  Neuen LaunchDaemon installieren (benötigt sudo):
    sudo cp com.fueltracker.plist /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/com.fueltracker.plist
    sudo chmod 644 /Library/LaunchDaemons/com.fueltracker.plist
    sudo launchctl bootstrap system /Library/LaunchDaemons/com.fueltracker.plist

  Status prüfen:  sudo launchctl print system/com.fueltracker | grep -E "state|exit"
  Stoppen:        sudo launchctl bootout system/com.fueltracker
  Starten:        sudo launchctl bootstrap system /Library/LaunchDaemons/com.fueltracker.plist
  Logs:           tail -f /Users/uwe/Projects/fuel-tracker/logs/stdout.log
-->
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.fueltracker</string>

  <!-- Als Benutzer uwe ausführen, nicht als root -->
  <key>UserName</key>
  <string>uwe</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/Caskroom/miniforge/base/envs/fuel-tracker/bin/python</string>
    <string>/Users/uwe/Projects/fuel-tracker/main.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/uwe/Projects/fuel-tracker</string>

  <!-- Beim Mac-Start automatisch laden -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Bei Absturz neu starten -->
  <key>KeepAlive</key>
  <true/>

  <!-- Log-Dateien -->
  <key>StandardOutPath</key>
  <string>/Users/uwe/Projects/fuel-tracker/logs/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/uwe/Projects/fuel-tracker/logs/stderr.log</string>

</dict>
</plist>
